# packaging/xoa-proxy.service
[Unit]
Description=XOA Deploy Proxy — HTTP/gunzip bridge for XAPI VM.import
Documentation=https://github.com/vagrantin/xoa-proxy
# Start after the network is up
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/opt/xensource/bin/xoa-proxy
Restart=on-failure
RestartSec=5s

# ── Configuration ────────────────────────────────────────────────────────
# Certificate verification is ON by default.
#
# If your XOA image is hosted behind a self-signed certificate, choose ONE of:
#
#Option A — CLI flag (edit ExecStart above):
#  ExecStart=/opt/xensource/bin/xoa-proxy --no-verify-ssl
#
#Option B — Environment variable
#  Environment=XOA_PROXY_VERIFY_SSL=0
#

# ── Logging ─────────────────────────────────────────────────────────────────
StandardOutput=null
StandardError=null
SyslogIdentifier=xoa-proxy

# ── Minimal privilege hardening ──────────────────────────────────────────────
# The proxy only needs outbound HTTPS and to listen on a loopback port > 1024.
# No filesystem writes, no setuid, no new privileges.
NoNewPrivileges=yes
PrivateTmp=yes

# Read-only access to the binary itself; no write access to the rest of the FS.
ProtectSystem=full
ProtectHome=yes

# Bind only to loopback.
RestrictAddressFamilies=AF_INET AF_INET6
AmbientCapabilities=

[Install]
WantedBy=multi-user.target
